# configs/systemd/user/empty-trash.service
[Unit]
Description=Empty user trash

[Service]
Type=oneshot
ExecStart=/usr/bin/gio trash --empty
User=jon
